# ansible/roles/common/templates/etc/cron.hourly/ntpdate.j2
#!/bin/bash
### managed by ansible

NTPDATE=$(which ntpdate)

$NTPDATE a.ntp.br
